# packaging/systemd/smcv.service
[Unit]
Description=SMCV encrypted secrets vault
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=smcv
Group=smcv
UMask=0077
EnvironmentFile=/etc/smcv/smcv.env
StateDirectory=smcv smcv-key
StateDirectoryMode=0700
ExecStartPre=/usr/local/lib/smcv/smcv-server preflight
ExecStart=/usr/local/lib/smcv/smcv-server
Restart=on-failure
RestartSec=5s
TimeoutStopSec=20s
LimitCORE=0
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
RemoveIPC=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProcSubset=pid
ProtectSystem=strict
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
CapabilityBoundingSet=
IPAddressDeny=any
IPAddressAllow=localhost

[Install]
WantedBy=multi-user.target
